# docker/supervisor.conf
[program:gemmolo]
command = gunicorn --bind 0.0.0.0:8000 gem.wsgi
environment = DJANGO_SETTINGS_MODULE="gem.settings.docker"
directory = /deploy/
redirect_stderr = true
